# etc/cron.d/perso
0 0 * * * nicolas ~/config/tools/db.sh > /var/log/cron.log 2>&1